# PI-4/Agente de monitorizacion/soc-sensor.service
[Unit]
# Nombre descriptivo del servicio que aparece en journalctl y systemctl status
Description=SOC Node Sensor - Agente de monitorización Pi4-Felix
# No arrancar hasta que la red esté disponible (necesitamos conexión con AWS)
After=network.target

[Service]
# Simple = el proceso principal ES el agente, systemd lo vigila directamente
Type=simple

# Ejecutar como root para poder leer /var/log/auth.log y usar iptables
User=root

# Directorio de trabajo: aquí debe estar config.yml y los certificados
WorkingDirectory=/home/lopex/pi4-felix/

# Usar el Python del entorno virtual, no el del sistema
ExecStart=/home/lopex/pi4-felix/env/bin/python /home/lopex/pi4-felix/agente_monitor3.py

# Reiniciar solo si el script falla con error (no si lo paras tú con systemctl stop)
Restart=on-failure

# Esperar 10 segundos antes de reiniciar (evita bucles de reinicios rápidos)
RestartSec=10

# Redirigir la salida estándar al journal (ver con: journalctl -u soc-sensor -f)
StandardOutput=journal
StandardError=journal

[Install]
# Arrancar automáticamente cuando el sistema esté en modo multiusuario (arranque normal)
WantedBy=multi-user.target